# coffee_watcher.service
[Unit]
Description=Coffee watcher service

[Service]
ExecStart=/home/pi/shipit-coffee-machine/start_watcher.sh

[Install]
WantedBy=multi-user.target